# teleclaw.service
[Unit]
Description=TeleClaw - Telegram remote control for Claude Code
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=%h/.local/bin/teleclaw
WorkingDirectory=%h
Restart=on-failure
RestartSec=10
Environment=PYTHONUNBUFFERED=1

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=teleclaw

[Install]
WantedBy=default.target
